# deploy/amazon-tracker-bot.service
[Unit]
Description=Amazon Tracker Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/amazon-tracker-bot
EnvironmentFile=/opt/amazon-tracker-bot/.env
ExecStart=/opt/amazon-tracker-bot/.venv/bin/python /opt/amazon-tracker-bot/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
